# supervisor.conf
[program:my_python_app]
command = /usr/bin/python /app/server.py
stdout_logfile = /app/server.log
stderr_logfile = /app/server-err.log
